# inference/api/deploy/systemd/vllm.service
[Unit]
Description=KrishiVaidya vLLM backend (GPU-profiled)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/krishivaidya/model
ExecStart=/opt/krishivaidya/model/.venv/bin/python -m inference.api.scripts.launch_vllm --port 8091
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
